# tools/stunnel.service.in
[Unit]
Description=TLS tunnel for network daemons
After=syslog.target network-online.target

[Service]
LimitNOFILE=20480
ExecStart=@bindir@/stunnel
ExecReload=/bin/kill -HUP $MAINPID
Type=forking
PrivateTmp=true

[Install]
WantedBy=multi-user.target
